# .config/systemd/user/open-rgb.service
[Unit]
Description=OpenRGB control

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb -m 'Color Cycle'

[Install]
WantedBy=default.target
